# deploy/vps/systemd/frenchnexora-api.service
[Unit]
Description=Nexora French movie API (nexora-node-api)
After=network.target

[Service]
Type=simple
User=codex
WorkingDirectory=/opt/nexora/frenchnexoraAPI
Environment=NODE_ENV=production
Environment=PORT=3100
Environment=PUPPETEER_SKIP_CHROMIUM_DOWNLOAD=true
Environment=PUPPETEER_EXECUTABLE_PATH=/usr/bin/chromium
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
