# deploy/vps/systemd/frenchnexora-api.service
[Unit]
Description=Nexora French streaming providers API
After=network.target

[Service]
Type=simple
User=nexora
WorkingDirectory=/opt/nexora/frenchnexoraAPI
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=API_PROVIDER_TIMEOUT_MS=55000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
